# alpine/etc/profile.d/color_prompt.sh.disabled
version https://git-lfs.github.com/spec/v1
oid sha256:ba24425c6864a5d17fa0fdaf914c4d21419e47c4d62080c33830af059fe46617
size 447
