# alpine/etc/profile.d/locale.sh
version https://git-lfs.github.com/spec/v1
oid sha256:84eb9034099d759ff08e6da5a731cacfc63a319547ad0f1dfc1c64853aca93f2
size 61
